# profile.d/Gulp.ps1
Invoke-Expression ((gulp --completion=powershell) -join [System.Environment]::NewLine)
